# crontab
*/1 * * * * /usr/local/bin/python /cronsql.py >> /proc/1/fd/1 2>/proc/1/fd/2
